# eventum-monitor.cron
# eventum monitor
# report status to root
MAILTO=root

*/10 * * * * http exec /usr/share/eventum/crons/monitor.php -q
