# eventum-monitor.cron
# eventum reminder
# bounce messages back to root
MAILTO=root

*/10 * * * * http cd /usr/share/eventum/misc && /usr/bin/php4 -q monitor.php

# vim:syn=crontab
