# eventum-monitor.cron
# eventum reminder
# bounce messages back to root
MAILTO=root

*/10 * * * * eventum cd /usr/share/eventum && /usr/bin/php4 -q monitor.php

# vim:syn=crontab
