# eventum-reminder.cron
# eventum reminder
# bounce messages back to root
MAILTO=root

*/10 * * * * http /usr/share/eventum/check_reminders.php

# It is recommended that you run the reminder cron job every 10 minutes,
# so it won't flood you with alerts too often, but it would still be
# enough to handle most cases.
