# eventum-monitor.cron
# eventum reminder
# bounce messages back to root
MAILTO=root

*/10 * * * * http /usr/share/eventum/monitor.php
