# eventum-monitor.cron
# eventum reminder
# bounce messages back to root
MAILTO=root

*/10 * * * * http exec /usr/share/eventum/crons/monitor.php
